# kalien-verifier/deploy/supervisord/kalien-api.conf
; NOTE: Update command and directory below to match your actual clone path.
; The VASTAI script defaults to /workspace/kalien, but a manual
; `git clone` may use a different name.
;
; Environment variables are sourced from /etc/kalien/api-server.env.
; Edit that file to set API_KEY and other config.

[program:kalien-api]
command=/bin/bash -c "set -a && source /etc/kalien/api-server.env && set +a && exec /workspace/kalien/kalien-verifier/target/release/api-server"
directory=/workspace/kalien/kalien-verifier
autostart=true
autorestart=true
startsecs=5
startretries=999
priority=100
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/lib/kalien/prover/api-server.log
stderr_logfile=/var/lib/kalien/prover/api-server.err
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
